# meshnode/overlay/etc/udev/rules.d/99-mesh-network.rules
ACTION=="add", DRIVERS=="ath9k_htc", SUBSYSTEM=="net", KERNEL=="wlan*", RUN+="/usr/bin/setup-mesh %k"
